# klipper/moonraker-telegram-bot.ini
[program:moonraker-telegram-bot]
user=klippy
command=/home/klippy/moonraker-telegram-bot-env/bin/python /home/klippy/moonraker-telegram-bot/bot/main.py -c /home/klippy/.config/telegram.conf
environment=USER=klippy,HOME=/home/klippy
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes = 0
autorestart=true
autostart=false
dependent_startup=true
dependent_startup_wait_for=moonraker:running
